# applications/luci-app-DeviceReboot/root/etc/uci-defaults/luci-reboot
#!/bin/sh

uci -q batch <<-EOF >/dev/null
        add ucitrack reboot
        set ucitrack.@reboot[-1].init=reboot
        commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
